# system/cron.daily/blocklist-de
#!/bin/sh

# This script refreshes the block lists in /etc/blocklist.de from blocklist.de
# The list of blocklists can be found on https://www.blocklist.de/en/export.html
#
# The IPs are used in /etc/hosts.deny. See man 5 hosts_allow for more details

wget -O /etc/blocklist.de/strongips.txt https://lists.blocklist.de/lists/strongips.txt
wget -O /etc/blocklist.de/ssh.txt https://lists.blocklist.de/lists/ssh.txt
